# src/stack_gps/tools/base_station/99-ublox-f9p.rules
# u-blox ZED-F9P 자동 바인딩 규칙
# 이 보드는 커스텀 USB ID(1546:0507 / 1546:0508)를 써서 ftdi_sio가 자동 인식하지 못한다.
# USB 장치가 꽂히면 ftdi_sio를 로드하고 해당 ID를 new_id에 등록하여 ttyUSB 포트를 생성한다.

# --- USB 장치 감지 시: ftdi_sio 로드 + 커스텀 ID 등록 ---
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="1546", ATTR{idProduct}=="0507", RUN+="/sbin/modprobe ftdi_sio", RUN+="/bin/sh -c '/bin/echo 1546 0507 > /sys/bus/usb-serial/drivers/ftdi_sio/new_id'"
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="1546", ATTR{idProduct}=="0508", RUN+="/sbin/modprobe ftdi_sio", RUN+="/bin/sh -c '/bin/echo 1546 0508 > /sys/bus/usb-serial/drivers/ftdi_sio/new_id'"

# --- ModemManager 차단 ---
# USB 재연결 시 ModemManager가 모뎀인지 조사하느라 포트를 수십 초간 점유(EBUSY)한다.
# u-blox GNSS는 모뎀이 아니므로 조사 대상에서 제외.
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="1546", ENV{ID_MM_DEVICE_IGNORE}="1"

# --- ttyUSB 포트 생성 시: 권한 666 + 안정적 심볼릭 링크 ---
# UART1 (0507) = GPS 메인 통신 포트 → /dev/ttyF9P 로 고정
SUBSYSTEM=="tty", ATTRS{idVendor}=="1546", ATTRS{idProduct}=="0507", MODE="0666", GROUP="dialout", ENV{ID_MM_DEVICE_IGNORE}="1", SYMLINK+="ttyF9P"
# UART2 (0508)
SUBSYSTEM=="tty", ATTRS{idVendor}=="1546", ATTRS{idProduct}=="0508", MODE="0666", GROUP="dialout", ENV{ID_MM_DEVICE_IGNORE}="1", SYMLINK+="ttyF9P_uart2"
